# roles/glances/templates/glances.service.j2
[Unit]
Description = Glances in Web Server Mode
After = network.target

[Service]
ExecStart = {{ glances_bin[ansible_os_family] }} -w

[Install]
WantedBy = multi-user.target
